# flips.service
[Unit]
Description=Flips FastAPI Application
After=network.target

[Service]
Type=simple
# Этот файл будет автоматически создан скриптом setup-service.sh
# с правильными путями для вашей системы
WorkingDirectory=/path/to/flips/backend
ExecStart=/path/to/flips/backend/start.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
