# flips.service
[Unit]
Description=Flips FastAPI Application
After=network.target

[Service]
Type=simple
# Уберите комментарий и укажите правильного пользователя, если нужно запускать не от root
# User=your_username
WorkingDirectory=/root/flips/backend
# Используем скрипт-обертку для автоматического поиска uv
ExecStart=/root/flips/backend/start.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
