# flips.service
[Unit]
Description=Flips FastAPI Application
After=network.target

[Service]
Type=simple
# Уберите комментарий и укажите правильного пользователя, если нужно запускать не от root
# User=your_username
WorkingDirectory=/root/flips/backend
Environment="PATH=/root/.local/bin:/usr/local/bin:/usr/bin:/bin"
# Если uv установлен в другом месте, найдите его: which uv
ExecStart=/usr/local/bin/uv run uvicorn main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
